# claude-memory.service
[Unit]
Description=Claude Memory MCP Server
After=network.target

[Service]
Type=simple
User=dennis
WorkingDirectory=/home/dennis/claude-memory
EnvironmentFile=/home/dennis/claude-memory/.env
ExecStart=/home/dennis/.local/bin/uv run uvicorn server:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
